# touchpad-bt-toggle.service
[Unit]
Description=Toggle touchpad on Bluetooth mouse connection
After=bluetooth.target

[Service]
ExecStart=%h/.local/bin/touchpad-bt-toggle.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
